# systemd/m5/system/litellm-carmenta.service.d/memory.conf
# Brokkr · m5 memory limit — litellm-carmenta (LiteLLM proxy). Part of the Carmenta
# workshop stack; small (~90 MB observed) but bounded for safety. Applies live on
# daemon-reload (currently active, well under the cap → safe, no reclaim).
[Service]
MemoryHigh=1G
MemoryMax=2G
OOMScoreAdjust=400
