# systemd/m5/system/whisper-server.service.d/memory.conf
# Brokkr · m5 — whisper-server (speech inference). Important inference service:
# bounded (catches a runaway) but protected from the global kill list.
[Service]
MemoryHigh=3G
MemoryMax=4G
OOMScoreAdjust=-100
ManagedOOMPreference=avoid
