# systemd/m5/system/llama-swap.service.d/memory.conf
# Brokkr · m5 — PROTECT llama-swap (the model server; ~20 GB observed, more when a
# larger model is loaded). This is the box's whole reason to exist, so it must be the
# LAST thing the kernel kills, never the first.
#
# NO MemoryMax here on purpose — a hard cap would OOM-kill the model mid-inference when
# it legitimately grows to load a bigger model. Protection is via OOM ordering instead.
#
# MemoryMin reserves memory that won't be reclaimed under pressure. NOTE (cgroup v2):
# a child's effective memory.min is capped by its ancestors' — system.slice defaults to
# 0, so for this to bite you must ALSO protect the slice:
#   sudo systemctl set-property system.slice MemoryMin=20G
# The deploy script does this. Without it, MemoryMin here is a no-op (OOMScoreAdjust and
# the non-critical caps still do the real work).
[Service]
MemoryMin=20G
OOMScoreAdjust=-700
ManagedOOMPreference=avoid
