# systemd/m5/system/home-gateway.service.d/memory.conf
# Brokkr · m5 — home-gateway (the :8080 homeserver gateway fronting the model API).
# Critical for availability but small (~35 MB). A generous ceiling catches a leak
# without threatening it, and negative OOMScoreAdjust keeps it off the kill list.
[Service]
MemoryMax=2G
OOMScoreAdjust=-300
ManagedOOMPreference=avoid
