# systemd/orin/jtop.service.d/10-brokkr-cpu-limit.conf
# Brokkr resource guard for jetson-stats' jtop.service on orin (brokkr#15).
#
# jetson-stats 4.3.2 (the latest release on the PyPI path this host installs from
# — upgrading via pip is NOT a fix; upstream GitHub has newer 7.x if the host can
# consume it) has a known
# bug where the jtop daemon spins on a core. On 2026-07-02 it pegged a core for
# ~14h, drove load avg to ~21, and starved host services — notably wedging the
# Heimdall fleet push agent (heimdall-agent.service), which stayed
# active/running but stopped delivering pushes, showing orin as falsely
# "Offline" in Heimdall (brokkr#14 is the complementary watchdog).
#
# This drop-in is the substrate-layer guard (same philosophy as the m5
# per-service memory limits, brokkr#5): even if jtop spins again, the kernel
# caps it so it can no longer consume a whole core or starve neighbours.
#
# CPUQuota=40%   → jtop may use at most ~0.4 of one core (plenty to sample
#                  tegrastats; a spin can no longer peg a full core).
# CPUWeight=20   → under contention, real workloads and the push agent win.
# MemoryMax=256M → backstop; jtop's normal RSS is ~75M.
#
# Install (needs root on orin — jtop.service is a system-manager unit):
#   sudo install -Dm644 10-brokkr-cpu-limit.conf \
#     /etc/systemd/system/jtop.service.d/10-brokkr-cpu-limit.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart jtop.service   # also clears any current spin
# Verify:  systemctl show jtop.service -p CPUQuotaPerSecUSec -p MemoryMax
[Service]
CPUQuota=40%
CPUWeight=20
MemoryMax=256M
