# systemd/m5/system/ttyd-carmenta.service.d/memory.conf
# Brokkr · m5 memory limit — ttyd-carmenta (web terminal). NON-CRITICAL, tiny.
# Bounded and killed early on global OOM. Applies on next start.
[Service]
MemoryHigh=512M
MemoryMax=1G
# Kernel OOM ordering — dies early on a global OOM. No ManagedOOMPreference (see the
# workshop-carmenta drop-in for why `omit` would be backwards here).
OOMScoreAdjust=500
