# systemd/m5/user/workshop-carmenta.service.d/memory.conf
# Brokkr · m5 memory limit — workshop-carmenta (open-webui). NON-CRITICAL.
# Runaway suspect in the 2026-07-01 OOM storm (brokkr#5). Hard-bounded so it can
# never starve llama-swap or an unattended overnight run; killed FIRST on global OOM.
# Applies on next start (service is normally inactive). Tune from MemoryPeak — see
# docs/m5-memory-limits.md.
[Service]
MemoryHigh=4G
MemoryMax=6G
# Kernel OOM ordering (this is the backstop that makes it die first on a global OOM).
OOMScoreAdjust=600
# NOTE: no ManagedOOMPreference here. In systemd-oomd, `omit` EXCLUDES a cgroup from kill
# candidates (i.e. protects it) — the opposite of what we want for a first-to-die service;
# the default (`none`) leaves it a normal candidate. Kill ordering is driven by
# OOMScoreAdjust above, which works regardless of whether systemd-oomd is enabled.
