# systemd/m5/user/code-carmenta.service.d/memory.conf
# Brokkr · m5 memory limit — code-carmenta (workshop code-server + its vitest runs).
# NON-CRITICAL. Runaway suspect in the 2026-07-01 OOM storm (brokkr#5). vitest and other
# dev processes spawned by code-server inherit THIS cgroup, so this cap bounds them too.
# Killed first on global OOM; applies on next start. Tune from MemoryPeak.
[Service]
MemoryHigh=3G
MemoryMax=5G
# Kernel OOM ordering — dies first on a global OOM. No ManagedOOMPreference: `omit` would
# PROTECT it from systemd-oomd (the inverse of what we want); default leaves it a candidate.
OOMScoreAdjust=700
